# crontab
MAILTO=""
0 9 * * * zsh -ic "update--all" > /dev/null
